# whatspricebot.service.example
[Unit]
Description=WhatsPriceBot WhatsApp Price Bot
After=network.target

[Service]
WorkingDirectory=/opt/whatspricebot
EnvironmentFile=/opt/whatspricebot/.env
ExecStart=/opt/whatspricebot/venv/bin/python /opt/whatspricebot/app.py
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
